# scripts/com.claude-wormhole.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-wormhole.web</string>

    <!-- Update these paths to match your setup -->
    <key>WorkingDirectory</key>
    <string>/Users/YOU/www/claude-wormhole</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Use: which node -->
        <string>/Users/YOU/.nvm/versions/node/v24.13.0/bin/node</string>
        <string>dist/server.cjs</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <!-- Use: echo $PATH (from a terminal with node available) -->
        <key>PATH</key>
        <string>/Users/YOU/.nvm/versions/node/v24.13.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-wormhole.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-wormhole.err</string>
</dict>
</plist>
