# scripts/com.claude-wormhole.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-wormhole.web</string>

    <key>WorkingDirectory</key>
    <string>/Users/shyam/www/claude-wormhole</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/shyam/www/claude-wormhole/bin/wormhole</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>TZ</key>
        <string>Asia/Calcutta</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-wormhole.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-wormhole.err</string>
</dict>
</plist>
